# udev/01-detect-usb-key.rules
# DETECT USB STICK PLUG
ACTION=="add",KERNEL="sd[a-z]",RUN+="/sabu/usb-detect.sh"
